# symlinks/etc/systemd/system/netspeed.timer
[Unit]
Description=test

[Timer]
OnUnitActiveSec=5s
Persistent=true
AccuracySec=1us

[Install]
WantedBy=timers.target
